# deploy/systemd/units/drift-healthcheck.service
[Unit]
Description=Drift Stack Health Check

[Service]
Type=oneshot
ExecStart=/opt/drift-stack-systemd/bin/drift-healthcheck.sh
