# deploy/systemd/units/drift-panel.service
[Unit]
Description=Drift Experience Panel Static Server - 7x24
After=network-online.target
Wants=network-online.target
PartOf=drift-stack.target
StartLimitIntervalSec=300

[Service]
Type=simple
User=__STACK_USER__
EnvironmentFile=/etc/drift-stack.env
WorkingDirectory=/opt/drift-stack-systemd
ExecStart=/opt/drift-stack-systemd/bin/run-drift-panel.sh
Restart=always
RestartSec=3
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
WantedBy=drift-stack.target
